# worker/vigil-worker.service
[Unit]
Description=Vigil background worker (BullMQ)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/vigil
EnvironmentFile=/home/ubuntu/vigil/.env.local
ExecStart=/usr/bin/npm run worker
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
